# systemd/user/run-u411.service
# This is a transient unit file, created programmatically via the systemd API. Do not edit.
[Unit]
Description=/home/kate/.config/i3/timed-wallpaper.sh

[Service]
ExecStart="/home/kate/.config/i3/timed-wallpaper.sh"
Environment=DISPLAY=:0 
